# bwd.service
[Unit]
Description=Bing wallpaper downloader
Requisite=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bwd

[Install]
